# etc/systemd/system/nft@.service
[Unit]
Description=nft
Documentation=man:ntf(8)

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=/usr/sbin/nft -f /etc/firewall/%I.nft 
ExecStop=/usr/sbin/nft flush ruleset

[Install]
WantedBy=multi-user.target
DefaultInstance=filter
